# Services/kasmvnc.service
[Unit]
Description=KasmVNC Server for QuadRF Desktop
After=network.target

[Service]
Type=forking
User=dietpi
Group=dietpi
WorkingDirectory=/home/dietpi

# Clean up any existing lock files in case of a hard crash
ExecStartPre=-/usr/bin/vncserver -kill :1 > /dev/null 2>&1
ExecStartPre=-/bin/rm -f /tmp/.X1-lock
ExecStartPre=-/bin/rm -f /tmp/.X11-unix/X1

# Start the KasmVNC server
ExecStart=/usr/bin/vncserver :1 -disableBasicAuth

# Stop the KasmVNC server cleanly
ExecStop=/usr/bin/vncserver -kill :1

# Automatically restart if it crashes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
